# roles/samba-dc/templates/supervisor-samba.conf
[program:samba]
command=/usr/local/samba/sbin/samba --foreground --configfile={{samba_install_dir}}/etc/smb.conf --model=prefork --option 'prefork children = {{ansible_processor_vcpus * 2}}'
