# roles/samba-dc/templates/supervisor-samba.conf
[program:samba]
command=/usr/local/samba/sbin/samba --foreground --model=prefork --option 'prefork children = {{ansible_processor_vcpus * 2}}'
